# publishing-templates.ini
[program:publishing-templates]
command=/var/publishing-templates/publishing-templates
directory=/var/publishing-templates
autostart=true
autorestart=true
startretries=3
